# roles/ipfs-cluster/templates/ipfs-cluster.service
[Unit]
Description=IPFS Cluster Service
After=network.target