# roles/ipfs-cluster/templates/ipfs-cluster.service
[Unit]
Description=IPFS Cluster Service
After=network.target

[Service]
LimitNOFILE={{ ipfs_cluster_fd_max }}